# templates/solr.ini.j2
[program:solr]
directory = /opt/solr/example
command = java -jar -Djetty.port={{ solr_port }} start.jar
autostart = {{ solr_autostart }}
